# cron/crontab
##
# We're going to load the environment variables into our intermediary script
#

* * * * * BASH_ENV=/var/www/html/.env /opt/scripts/runcron.sh

# Cron files must have an ending blank line
